# package/deb/postinstall
# Initial installation
  systemctl --no-reload preset pds-server.service &>/dev/null || :
