# package/rpm/postinstall
if [ $1 -eq 1 ] ; then
  # Initial installation
  systemctl --no-reload preset pds-service.service &>/dev/null || :
fi
