# package/pds-server.service.erb
[Unit]
Description=Puppet Data Service server
After=network.target

[Service]
Type=simple
User=pds-server
WorkingDirectory=/opt/puppetlabs/server/apps/pds-server
Environment=RUBYOPT=-r/opt/puppetlabs/server/apps/pds-server/vendor/bundle/bundler/setup.rb
Environment=RUBYGEMS_GEMDEPS=/opt/puppetlabs/server/apps/pds-server/Gemfile
Environment=PATH=/opt/puppetlabs/puppet/bin:$PATH
Environment=GEM_PATH=/opt/puppetlabs/server/apps/pds-server/vendor/bundle/ruby/<%= ruby_version %>/
ExecStart=/opt/puppetlabs/server/apps/pds-server/vendor/bundle/ruby/<%= ruby_version %>/bin/puma -C config/puma.rb -e production
Restart=always

[Install]
WantedBy=multi-user.target
