# package/pds-server.service.erb
[Unit]
Description=Puppet Data Service server
After=network.target

[Service]
Type=simple
User=pds-server
WorkingDirectory=/opt/puppetlabs/server/apps/pds-server
ExecStart=/opt/puppetlabs/puppet/bin/ruby -r ./vendor/bundle/bundler/setup.rb ./vendor/bundle/ruby/<%= ruby_version %>/gems/puma-<%= puma_version %>/bin/puma -C config/puma.rb -e production
Restart=always

[Install]
WantedBy=multi-user.target
